# deploy/quadlets/sports-log.container
[Unit]
Description=Sports Log API
After=postgres.service
Requires=postgres.service

[Container]
Image=localhost/sports-log:latest
ContainerName=sports-log

PublishPort=3000:3000

# App config — set these in /etc/sports-log/app.env
EnvironmentFile=/etc/sports-log/app.env

Network=sports-log-net.network

HealthCmd=curl -sf http://localhost:3000/metrics || exit 1
HealthInterval=15s
HealthTimeout=5s
HealthRetries=3

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target default.target
